# s6/services.d/supercronic/run
#!/usr/bin/with-contenv sh

CRONTAB="/etc/crontab"

"/usr/local/bin/supercronic" "$CRONTAB"
