# s6/services.d/supercronic/run
#!/usr/bin/with-contenv bash

CRONTAB="/etc/crontab"

exec \
  s6-setuidgid abc \
    "/usr/local/bin/supercronic" \
    "$CRONTAB"
